# debian/ppc-session-logger.user.service
[Unit]
Description=Pardus Parental Control - Session Logger Service
Requisite=graphical-session.target
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=exec
Restart=no
Slice=app.slice
ExitType=cgroup
ExecStart=/usr/bin/pkexec /usr/share/pardus/pardus-parental-control/src/session_logger.py %u

[Install]
WantedBy=graphical-session.target
